# etc/cron.d/myntcd
#Daemon kimenetenek feldolgozasa:
*/5 * * * * root /usr/local/myntcd/bin/myntc.php
#Szuro aktivalasa:
*/5 * * * * root /usr/local/myntcd/bin/blacklist.php
