# tcws/etc/cron.d/tcrules
2-57/5 * * * * root /usr/local/bin/tcrules.py
